# assets/flamingo.target
[Unit]
Description=Flamingo target allowing to start/stop all flamingo@.service instances at once
